# standalone-local-api/com.bizneai.localapi.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.bizneai.localapi</string>
  <key>WorkingDirectory</key>
  <string>/usr/local/lib/bizneai-local-api</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>BIZNEAI_EMBEDDED</key>
    <string>1</string>
    <key>BIZNEAI_USER_DATA</key>
    <string>/usr/local/lib/bizneai-local-api/data</string>
    <key>PORT</key>
    <string>3000</string>
  </dict>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>start.cjs</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/usr/local/lib/bizneai-local-api/logs/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/usr/local/lib/bizneai-local-api/logs/launchd.err.log</string>
</dict>
</plist>
